# example-use-case/client-example/alarm-repo-check.service
[Unit]
Description=Check git repo for upstream changes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alarm-repo-check.sh
